# debian/minios-welcome.path
[Unit]
Description=Watch browser configuration for changes

[Path]
PathChanged=/usr/share/firefox-esr/browser/defaults/preferences/firefox-branding.js
PathChanged=/etc/firefox/syspref.js

[Install]
WantedBy=multi-user.target
